# roles/multicat/templates/multicat_stream.service.j2
[Unit]
Description=Multicat streaming service
[Service]
Type=simple
ExecStart=/home/centos/multicat_loop.sh
[Install]
WantedBy=multi-user.target
